# etc/supervisord.conf
[supervisord]
nodaemon=true

[program:fhem]
command=/etc/init.d/fhem start

